# luci-app-gecoosac/root/etc/uci-defaults/luci-gecoosac
#!/bin/sh

PORT_HELPER=/usr/libexec/gecoosac-random-port

port="$(uci -q get gecoosac.config.port)"
changed=0

if [ -z "$port" ] || [ "$port" = 'auto' ]; then
	port="$("$PORT_HELPER")" || exit 1
	[ -n "$port" ] || exit 1
	uci -q set gecoosac.config.port="$port"
	changed=1
fi

[ "$changed" -eq 0 ] || uci -q commit gecoosac

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@gecoosac[-1]
	add ucitrack gecoosac
	set ucitrack.@gecoosac[-1].init=gecoosac
	commit ucitrack
EOF

exit 0
